# rc.local
#!/bin/sh

ln -s /data/dbus-sma-smartmeter/service /service/dbus-sma-smartmeter

exit 0
